# puma6/atom/mod/lib/systemd/system/ffinit.service
[Service]
ExecStart=/etc/init.d/S91-ffinit
Type=oneshot
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
