# atom/mod/lib/systemd/system/telnetd.service
[Service]
Type=forking
ExecStart=/usr/bin/telnetd
[Install]
WantedBy=multi-user.target
