# puma7/atom/mod/lib/systemd/system/ffstart.service
[Service]
ExecStart=/etc/init.d/S94-ffstart
Type=oneshot
After=ffimage.service
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
