# puma7/atom/mod/lib/systemd/system/ffusr.service
[Service]
ExecStart=/etc/init.d/S97-ffusr
Type=oneshot
After=ffstart.service
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
